# deploy/systemd/hashfarm-p2pool.service
[Unit]
Description=P2Pool v4 sidechain + stratum
After=hashfarm-monerod.service network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/Projects/HashFarm
EnvironmentFile=%h/Projects/HashFarm/.env
ExecStart=%h/Projects/HashFarm/scripts/garuda/p2pool-v4.sh run
Restart=on-failure
RestartSec=20

[Install]
WantedBy=default.target
